# root/etc/services.d/template/run
#!/usr/bin/with-contenv bash
echo OwoWhatsThis
exit